# local.mac.dakokuscript.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>local.mac.dakokuscript</string>
    <key>WorkingDirectory</key>
    <string>/Users/hiroki-shiotsuka/work/selenium/dakoku</string>
    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/python3</string>
      <string>taikin.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>10</integer>
        <key>Weekday</key>
        <integer>2</integer>
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>10</integer>
        <key>Weekday</key>
        <integer>3</integer>
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>10</integer>
        <key>Weekday</key>
        <integer>4</integer>
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>10</integer>
        <key>Weekday</key>
        <integer>5</integer>
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>10</integer>
      </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/hiroki-shiotsuka/work/selenium/dakoku/log/access.log</string>
  </dict>
</plist>
